# crontab
00 * * * * USER=hatter /home/hatter/.dotfiles/backup.sh > /dev/null
30 * * * * /home/hatter/.dotfiles/sync_backups.sh
