# crontab
*/10 * * * * /home/hatter/.dotfiles/backup.sh > /dev/null
05 * * * * /home/hatter/.dotfiles/sync_backups.sh
