# vm_scripts/vm_setup.service
[Unit]
Description=VM Provisioning Service
After=auditd.service

[Service]
ExecStart=/usr/local/bin/vm_setup

[Install]
WantedBy=network.target
